# templates/elasticsearch.init.d.erb
#!/bin/sh
# chkconfig: 2345 20 80
# description: Handle elasticsearch
#

# file managed by puppet

<%
   es10 = @version.start_with? '1.'
   bootstrapClass = "Elastic#{es10 ? 's' : 'S'}earch"
%>

if test -f /lib/lsb/init-functions; then
    . /lib/lsb/init-functions
fi
if test -f /etc/init.d/functions; then
    . /etc/init.d/functions
fi

SHUTDOWN_WAIT=30

# try to find 

export JAVA_HOME=<%= @javahome %>

if [ "$( whoami )" != "<%= @user %>" ]; then
    # SU="su <%= @user %> -c"
    SU="daemon --user <%= @user %>"
else
    SU=""
fi

PIDFILE=<%= @es_pidfile %>

start() {
  isrunning

  if [ "$?" = 0 ]; then
    echo "elasticsearch already running"
    return 0
  fi

  # Change directory to prevent path problems
  cd <%= @es_path %>
  # Source JVM parameters
  test -r bin/elasticsearch.in.sh && . bin/elasticsearch.in.sh 

  # Remove pidfile if still around
  test -f $PIDFILE && rm -f $PIDFILE

  COMMAND="<%= @es_path %>/bin/elasticsearch <%= '-d' if es10 %>"
  if [ "$SU" != "" ]; then
      umask 0002
      $SU $COMMAND > /dev/null
  else
      umask 0002
      $COMMAND >/dev/null
  fi
}

stop() {
  isrunning

  if [ "$?" = 1 ]; then
    echo "elasticsearch already stopped"
    rm -f $PIDFILE # remove pidfile if still around
    return 0
  fi

  echo -n "Waiting for elasticsearch to exit (${SHUTDOWN_WAIT} sec.): "

  count=0
  until [ "$pid" = "" ] || [ $count -gt $SHUTDOWN_WAIT ]; do
    COMMAND="/bin/kill $pid" 
    if [ "$SU" != "" ]; then
      $SU "$COMMAND" > /dev/null
    else
      $COMMAND >/dev/null
    fi
      
    findpid

    echo -n "."
    sleep 3
    count=$((count+3))
  done

  echo ""

  if [ "$count" -gt "$SHUTDOWN_WAIT" ]; then
    echo "Forcing elasticsearch to stop"
    /bin/kill -9 $pid && sleep 5
  fi

  # check if elasticsearch is still around, this will be our exit status
  ! isrunning
}

findpid() {
  pid=""
  pid=$(pgrep -U <%= @user %> -f "<%= @es_path %>.*org.elasticsearch.bootstrap.<%=bootstrapClass%>")

  # validate output of pgrep
  if ! [ "$pid" = "" ] && ! [ "$pid" -gt 0 ]; then
    log_failure_msg "Unable to determine if elasticsearch is running"
    exit 1
  fi
}

isrunning() {
  findpid

  if [ "$pid" = "" ]; then
    return 1
  elif [ "$pid" -gt 0 ]; then
    return 0
  fi
}

case "$1" in
  start)
    start
    RETVAL=$?

    if [ "$RETVAL" = 0 ]; then
      log_success_msg "Started elasticsearch"
    else
      log_failure_msg "Not able to start elasticsearch"
    fi
    ;;

  stop)
    stop
    RETVAL=$?

    if [ "$RETVAL" = 0 ]; then
      log_success_msg "Stopped elasticsearch"
    else
      log_failure_msg "Not able to stop elasticsearch"
    fi
  ;;

  restart)
    stop
    sleep 5
    start
    RETVAL=$?

    if [ "$RETVAL" = 0 ]; then
      log_success_msg "Restarted elasticsearch"
    else
      log_failure_msg "Not able to restart elasticsearch"
    fi
  ;;

  status)
    isrunning
    RETVAL=$?

    if [ "$RETVAL" = 0 ]; then
      echo "elasticsearch (pid $pid) is running..."
    else
      echo "elasticsearch is stopped"
      RETVAL=3
    fi
  ;;

  *)
    echo "Usage: $0 {start|stop|restart|status}."
  ;;

esac

exit $RETVAL
